# build/provision/roles/_volumes/files/ephemeral.conf
description "Mount EC2 ephemeral drives (if applicable)"

start on filesystem
script
    /usr/local/sbin/mount-ephemeral
end script
