# build/provision/roles/_registry/files/docker-registry.conf
description "Private Docker registry service"

start on started docker
respawn

script
    exec docker run -e SETTINGS_FLAVOR=s3 \
        -e AWS_BUCKET=$( user-data-value -k REGISTRY_BUCKET ) \
        -e STORAGE_PATH=$( user-data-value -k REGISTRY_PATH ) \
	-p 5000:5000 registry
end script

post-stop exec sleep 10
